# surf-ufm-pkeyd.service
[Unit]
Description=SURF UFM PKey Daemon
Documentation=file:///usr/local/bin/surf_ufm_pkeyd.py
After=network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
ExecStart=/usr/bin/surf_ufm_pkeyd.py monitor --ufm-user admin --etcd-user root --poll-interval 30 --no-ssl-verify
Restart=always
RestartSec=10
TimeoutStopSec=30

# Security settings (running as root)
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
MemoryMax=512M
TasksMax=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=surf-ufm-pkeyd

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=GODEBUG=""
Environment=GOTRACEBACK=none

[Install]
WantedBy=multi-user.target